# supervisorctl/bot.conf
[program:telegrambot]
directory=/home/ubuntu/cas_bot
command=/usr/bin/python3 pkg/app.py


autostart=true
autorestart=true
environment=HOME="/home/ubuntu/",USER="ubuntu"

stderr_Logfile=/home/ubuntu/cas_bot/logfile_err.log
stdout_Logfile=/home/ubuntu/cas_bot/logfile.log
